# config/supervisord.conf
[supervisord]
logfile=/Users/aatifquamre/masterbot/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/Users/aatifquamre/masterbot/config/supervisord.pid
nodaemon=false

[unix_http_server]
file=/Users/aatifquamre/masterbot/config/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///Users/aatifquamre/masterbot/config/supervisor.sock

[program:freqtrade_mean_reversion]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_mean_reversion.json --strategy MeanReversionV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/mean_reversion.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/mean_reversion.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/Users/aatifquamre/masterbot/logs/mean_reversion.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_trend_follow]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_trend_follow.json --strategy TrendFollowV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/trend_follow.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/trend_follow.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/Users/aatifquamre/masterbot/logs/trend_follow.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_scalp]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_scalp.json --strategy ScalpV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/scalp.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/scalp.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/Users/aatifquamre/masterbot/logs/scalp.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_swing]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_swing.json --strategy SwingV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/swing.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/swing.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/Users/aatifquamre/masterbot/logs/swing.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_daily]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_daily.json --strategy DailyTrendV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/daily.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/daily.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/Users/aatifquamre/masterbot/logs/daily.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:qnt_reply_listener]
command=/Users/aatifquamre/masterbot/venv/bin/python /Users/aatifquamre/masterbot/qnt/memory/reply_listener.py
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=true
stderr_logfile=/Users/aatifquamre/masterbot/logs/qnt_reply_listener.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
stdout_logfile=/Users/aatifquamre/masterbot/logs/qnt_reply_listener.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
user=aatifquamre

[program:nats_subscriber]
command=/Users/aatifquamre/masterbot/venv/bin/python /Users/aatifquamre/masterbot/qnt/nats_subscriber.py
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=true
stderr_logfile=/Users/aatifquamre/masterbot/logs/nats_subscriber.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
stdout_logfile=/Users/aatifquamre/masterbot/logs/nats_subscriber.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
user=aatifquamre
environment=HOME="/Users/aatifquamre"

[program:self_healer]
command=/Users/aatifquamre/masterbot/automation/run_self_healer.sh
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=true
stdout_logfile=/Users/aatifquamre/masterbot/logs/self_healer_cron.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/Users/aatifquamre/masterbot/logs/self_healer_cron.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
user=aatifquamre

[program:freqtrade_micro]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_paper.json --config /Users/aatifquamre/masterbot/config/port_8085.json --strategy MicroScalpV1 --db-url sqlite:////Users/aatifquamre/masterbot/user_data/tradesv3_micro.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=true
stdout_logfile=/Users/aatifquamre/masterbot/logs/freqtrade_micro.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/Users/aatifquamre/masterbot/logs/freqtrade_micro_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:webui_server]
command=/Users/aatifquamre/masterbot/venv/bin/python /Users/aatifquamre/masterbot/webui/server.py
directory=/Users/aatifquamre/masterbot/webui
autostart=true
autorestart=true
stdout_logfile=/Users/aatifquamre/masterbot/logs/webui_server.stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/Users/aatifquamre/masterbot/logs/webui_server.stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
user=aatifquamre
environment=HOME="/Users/aatifquamre"
