# config/supervisord.conf
[supervisord]
logfile=/Users/aatifquamre/masterbot/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/Users/aatifquamre/masterbot/config/supervisord.pid
nodaemon=false

[unix_http_server]
file=/Users/aatifquamre/masterbot/config/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///Users/aatifquamre/masterbot/config/supervisor.sock

[program:freqtrade]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_paper.json --strategy SampleStrategy
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
startretries=3
stdout_logfile=/Users/aatifquamre/masterbot/logs/freqtrade.stdout.log
stderr_logfile=/Users/aatifquamre/masterbot/logs/freqtrade.stderr.log
user=aatifquamre

[program:qnt_reply_listener]
command=/Users/aatifquamre/masterbot/venv/bin/python /Users/aatifquamre/masterbot/qnt/memory/reply_listener.py
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=true
stderr_logfile=/Users/aatifquamre/masterbot/logs/qnt_reply_listener.log
stdout_logfile=/Users/aatifquamre/masterbot/logs/qnt_reply_listener.log
user=aatifquamre
