# config/supervisord.conf
[supervisord]
logfile=%(here)s/../logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=%(here)s/../config/supervisord.pid
nodaemon=false

[unix_http_server]
file=/tmp/masterbot_supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/masterbot_supervisord.sock

[program:freqtrade_mean_reversion]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_paper.json --config %(here)s/../config/config_mean_reversion.json --strategy MeanReversionV1 --strategy-path %(here)s/../strategies/active --db-url postgresql://aatifquamre:dummy@localhost/masterbot_mean_reversion
directory=%(here)s/..
autostart=true
autorestart=unexpected
stdout_logfile=%(here)s/../logs/mean_reversion.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/mean_reversion.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_trend_follow]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_paper.json --config %(here)s/../config/config_trend_follow.json --strategy TrendFollowV1 --strategy-path %(here)s/../strategies/active --db-url postgresql://aatifquamre:dummy@localhost/masterbot_trend_follow
directory=%(here)s/..
autostart=true
autorestart=unexpected
stdout_logfile=%(here)s/../logs/trend_follow.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/trend_follow.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_scalp]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_paper.json --config %(here)s/../config/config_scalp.json --strategy ScalpV1 --strategy-path %(here)s/../strategies/active --db-url postgresql://aatifquamre:dummy@localhost/masterbot_scalp
directory=%(here)s/..
autostart=true
autorestart=unexpected
stdout_logfile=%(here)s/../logs/scalp.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/scalp.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_swing]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_paper.json --config %(here)s/../config/config_swing.json --strategy SwingV1 --strategy-path %(here)s/../strategies/active --db-url postgresql://aatifquamre:dummy@localhost/masterbot_swing
directory=%(here)s/..
autostart=true
autorestart=unexpected
stdout_logfile=%(here)s/../logs/swing.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/swing.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_daily]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_paper.json --config %(here)s/../config/config_daily.json --strategy DailyTrendV1 --strategy-path %(here)s/../strategies/active --db-url postgresql://aatifquamre:dummy@localhost/masterbot_daily
directory=%(here)s/..
autostart=true
autorestart=unexpected
stopwaitsecs=15
stdout_logfile=%(here)s/../logs/daily.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/daily.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:qnt_reply_listener]
command=%(here)s/../venv/bin/python %(here)s/../qnt/memory/reply_listener.py
directory=%(here)s/..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../logs/qnt_reply_listener.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
stdout_logfile=%(here)s/../logs/qnt_reply_listener.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
user=aatifquamre

[program:nats_subscriber]
command=%(here)s/../venv/bin/python %(here)s/../qnt/nats_subscriber.py
directory=%(here)s/..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../logs/nats_subscriber.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
stdout_logfile=%(here)s/../logs/nats_subscriber.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
user=aatifquamre
environment=HOME="/Users/aatifquamre"

[program:self_healer]
command=%(here)s/../automation/run_self_healer.sh
directory=%(here)s/..
autostart=true
autorestart=true
stdout_logfile=%(here)s/../logs/self_healer_cron.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=%(here)s/../logs/self_healer_cron.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
user=aatifquamre

[program:freqtrade_micro]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_paper.json --config %(here)s/../config/config_micro.json --strategy MicroScalpV1 --db-url postgresql://aatifquamre:dummy@localhost/masterbot_micro
directory=%(here)s/..
autostart=true
autorestart=true
stdout_logfile=%(here)s/../logs/freqtrade_micro.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/freqtrade_micro_err.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_bear_scalp]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_paper.json --config %(here)s/../config/config_bear_scalp.json --strategy BearScalpV1 --strategy-path %(here)s/../strategies/active --db-url postgresql://aatifquamre:dummy@localhost/masterbot_bear_scalp
directory=%(here)s/..
autostart=false
autorestart=unexpected
stdout_logfile=%(here)s/../logs/bear_scalp.stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=%(here)s/../logs/bear_scalp.stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
user=aatifquamre

[program:freqtrade_hyperliquid]
command=%(here)s/../venv/bin/freqtrade trade --config %(here)s/../config/config_hyperliquid.json --strategy MeanReversionV1 --strategy-path %(here)s/../strategies/active --db-url postgresql://aatifquamre:dummy@localhost/masterbot_hyperliquid
directory=%(here)s/..
autostart=false
autorestart=true
stdout_logfile=%(here)s/../logs/hyperliquid.stdout.log
stderr_logfile=%(here)s/../logs/hyperliquid.stderr.log
user=aatifquamre

[program:freqtrade_exporter]
command=%(here)s/../venv/bin/python -u %(here)s/../qnt/cockpit/exporter.py
directory=%(here)s/..
autostart=true
autorestart=true
stdout_logfile=%(here)s/../logs/exporter.stdout.log
stderr_logfile=%(here)s/../logs/exporter.stderr.log
user=aatifquamre

[program:prometheus]
command=/opt/homebrew/bin/prometheus --config.file=%(here)s/../config/prometheus.yml --storage.tsdb.path=%(here)s/../user_data/prometheus_data
directory=%(here)s/..
autostart=true
autorestart=true
stdout_logfile=%(here)s/../logs/prometheus.stdout.log
stderr_logfile=%(here)s/../logs/prometheus.stderr.log
user=aatifquamre

[program:obi_pipeline]
command=%(here)s/../venv/bin/python -u %(here)s/../qnt/oracle/realtime_orderflow.py
directory=%(here)s/..
autostart=true
autorestart=true
stdout_logfile=%(here)s/../logs/obi_pipeline.stdout.log
stderr_logfile=%(here)s/../logs/obi_pipeline.stderr.log
user=aatifquamre

[program:grafana]
command=/opt/homebrew/bin/grafana server --config /opt/homebrew/etc/grafana/grafana.ini --homepath /opt/homebrew/opt/grafana/share/grafana
directory=%(here)s/..
autostart=false
autorestart=true
stdout_logfile=%(here)s/../logs/grafana.stdout.log
stderr_logfile=%(here)s/../logs/grafana.stderr.log
user=aatifquamre

[program:telegram_gateway]
command=%(here)s/../venv/bin/python %(here)s/../qnt/memory/telegram_webhook_server.py --polling --port 8443
directory=%(here)s/..
autostart=true
autorestart=true
stderr_logfile=%(here)s/../logs/telegram_gateway.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
stdout_logfile=%(here)s/../logs/telegram_gateway.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
user=aatifquamre

[program:shadow_hyperopt]
command=%(here)s/../venv/bin/python %(here)s/../qnt/shadow/shadow_hyperopt.py
directory=%(here)s/..
autostart=true
autorestart=true
stdout_logfile=%(here)s/../logs/shadow_hyperopt.stdout.log
stderr_logfile=%(here)s/../logs/shadow_hyperopt.stderr.log
user=aatifquamre
