# config/supervisord.conf
[supervisord]
logfile=/Users/aatifquamre/masterbot/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/Users/aatifquamre/masterbot/config/supervisord.pid
nodaemon=false

[unix_http_server]
file=/Users/aatifquamre/masterbot/config/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///Users/aatifquamre/masterbot/config/supervisor.sock

[program:freqtrade_mean_reversion]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_mean_reversion.json --strategy MeanReversionV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/mean_reversion.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/mean_reversion.stdout.log
stderr_logfile=/Users/aatifquamre/masterbot/logs/mean_reversion.stderr.log
user=aatifquamre

[program:freqtrade_trend_follow]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_trend_follow.json --strategy TrendFollowV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/trend_follow.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/trend_follow.stdout.log
stderr_logfile=/Users/aatifquamre/masterbot/logs/trend_follow.stderr.log
user=aatifquamre

[program:freqtrade_scalp]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_scalp.json --strategy ScalpV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/scalp.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/scalp.stdout.log
stderr_logfile=/Users/aatifquamre/masterbot/logs/scalp.stderr.log
user=aatifquamre

[program:freqtrade_swing]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_swing.json --strategy SwingV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/swing.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/swing.stdout.log
stderr_logfile=/Users/aatifquamre/masterbot/logs/swing.stderr.log
user=aatifquamre

[program:freqtrade_daily]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_daily.json --strategy DailyTrendV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/daily.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/daily.stdout.log
stderr_logfile=/Users/aatifquamre/masterbot/logs/daily.stderr.log
user=aatifquamre

[program:qnt_reply_listener]
command=/Users/aatifquamre/masterbot/venv/bin/python /Users/aatifquamre/masterbot/qnt/memory/reply_listener.py
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=true
stderr_logfile=/Users/aatifquamre/masterbot/logs/qnt_reply_listener.log
stdout_logfile=/Users/aatifquamre/masterbot/logs/qnt_reply_listener.log
user=aatifquamre

[program:freqtrade_micro]
command=/Users/aatifquamre/masterbot/venv/bin/freqtrade trade --config /Users/aatifquamre/masterbot/config/config_micro.json --strategy MicroScalpV1 --strategy-path /Users/aatifquamre/masterbot/strategies/active --db-url sqlite:///user_data/micro.sqlite
directory=/Users/aatifquamre/masterbot
autostart=true
autorestart=unexpected
stdout_logfile=/Users/aatifquamre/masterbot/logs/micro.stdout.log
stderr_logfile=/Users/aatifquamre/masterbot/logs/micro.stderr.log
user=aatifquamre
